# laradock/php-worker/websockets.conf
[program:websockets]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/toko-kelontong websockets:serve
autostart=true
autorestart=true
numprocs=1
user=laradock
redirect_stderr=true